# cron/conf/cron/crontab
* * * * * cd /var/www/bearpass && php artisan schedule:run >> /dev/null 2>&1